# services/logging/usb-log-sync-periodic.service
[Unit]
Description=Periodic sync of recent openTPT logs to USB
# Only run if USB is mounted
ConditionPathIsMountPoint=/mnt/usb

[Service]
Type=oneshot
# Sync only recent logs (last 2 hours) to avoid large files
ExecStart=/usr/local/bin/usb-log-sync.sh
TimeoutStartSec=30
User=root
