# services/logging/usb-log-sync.timer
[Unit]
Description=Periodic USB log sync for openTPT
# Only run if USB is mounted
ConditionPathIsMountPoint=/mnt/usb

[Timer]
# Sync every 30 seconds for near real-time logging
OnBootSec=30s
OnUnitActiveSec=30s
Unit=usb-log-sync-periodic.service
AccuracySec=5s

[Install]
WantedBy=timers.target
