# config/boot/splash.service
[Unit]
Description=openTPT Boot Splash
DefaultDependencies=no
After=sysinit.target
Before=openTPT.service

[Service]
Type=oneshot
RemainAfterExit=yes
# Wait up to 5 seconds for framebuffer to become available
ExecStart=/bin/bash -c 'for i in $(seq 1 50); do [ -e /dev/fb0 ] && break; sleep 0.1; done; exec /usr/bin/fbi -T 1 -d /dev/fb0 --noverbose -a /home/pi/open-TPT/assets/splash.png'
StandardOutput=null
StandardError=null

[Install]
WantedBy=sysinit.target
